# .config/upstart/configTouchpad.conf
description "Configure touchpad"
start on desktop-start
task
script
    $HOME/.local/bin/cfgtouchpad.sh
end script
